# autostart.sh
#!/bin/sh
#chkconfig: 2345 80 90
#description:auto_run

cd /home/jzx/sample/WebRTC
nohup npm start&

